# .conf/supervisor-app.conf
[program:uwsgi]
command = uwsgi --ini /etc/uwsgi/sites/app.ini

[program:postgresql]
command = services start postgresql

[program:nginx]
command = nginx